# data/init/xsplashaa-phase3.in
#! /bin/sh
#
### BEGIN INIT INFO
# Provides:          xsplashaa-phase3
# Required-Start:    $local_fs
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Set Xsplashaa Phase 3
### END INIT INFO

if [ -r "@initdir@/xsplashaa-functions" ]; then
    . "@initdir@/xsplashaa-functions"
else
    echo "Could not read file @initdir@/xsplashaa-functions."
    exit 1
fi

NAME=xsplashaa-phase3
DESC="Xorg graphical boot Phase3"

set -e

case "$1" in
  start)
	$XSPLASHAA_CLIENT --phase=3 --socket $XSPLASHAA_SOCKET
	;;
  stop)
	;;
  *)
	N=/etc/init.d/$NAME
	# echo "Usage: $N {start|stop|restart|reload|force-reload}" >&2
	echo "Usage: $N {start|stop}" >&2
	exit 1
	;;
esac

exit 0

